# crysgarage-admin-sync.service
[Unit]
Description=CrysGarage Admin Firebase User Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/var/www/crysgarage-admin/backend
EnvironmentFile=/var/www/crysgarage-admin/backend/.env
Environment=PYTHONPATH=/var/www/crysgarage-admin/backend
ExecStart=/var/www/crysgarage-admin/backend/.venv/bin/python /var/www/crysgarage-admin/backend/run_sync.py
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
