# crysgarage-analytics.service
[Unit]
Description=CrysGarage Analytics Service
After=network.target

[Service]
User=nginx
Group=nginx
WorkingDirectory=/var/www/crysgarage.studio/analytics-service
Environment="PATH=/usr/bin:/usr/local/bin"
Environment="PYTHONPATH=/var/www/crysgarage.studio"
Environment="DB_PATH=/var/www/crysgarage.studio/data/crysgarage.db"
ExecStart=/usr/bin/python3 -m uvicorn analytics_app:app --host 127.0.0.1 --port 8005 --log-level info
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

